# roles/ferria_agent/templates/ferria-agent.timer.j2
# Managed by Ferria — manual changes will be overwritten on next convergence run.
# Generated by: roles/ferria_agent/templates/ferria-agent.timer.j2

[Unit]
Description=Ferria GitOps agent — periodic reconciliation timer
Documentation=https://github.com/The-Sec/ferria

[Timer]
# Fire 2 minutes after boot, then every N minutes thereafter.
OnBootSec=2min
OnUnitActiveSec={{ ferria_pull_interval_min }}min
{% if ferria_jitter_enabled | bool %}
# Thundering herd prevention: add a random delay before each run.
RandomizedDelaySec={{ ferria_jitter_max_sec }}
{% endif %}
# Catch up on missed runs after downtime (e.g., machine was off during a tick).
Persistent=true
# Allow up to 1 minute of timer drift to reduce wake-up pressure.
AccuracySec=1min

[Install]
WantedBy=timers.target
