# roles/ferria_agent/templates/ferria-agent.service.j2
# Managed by Ferria — manual changes will be overwritten on next convergence run.
# Generated by: roles/ferria_agent/templates/ferria-agent.service.j2

[Unit]
Description=Ferria GitOps agent — convergence run
After=network-online.target
Wants=network-online.target
Documentation=https://github.com/The-Sec/ferria

[Service]
Type=oneshot
ExecStart=/opt/ferria/ferria-wrapper.sh

# Give convergence runs up to 30 minutes before killing.
# Increase this if your playbooks are known to take longer.
TimeoutStartSec=1800

# Security hardening.
# Note: ProtectSystem is intentionally omitted — ansible-pull must write
# to /etc, /usr, and /var to install packages and manage configuration.
PrivateTmp=true
NoNewPrivileges=true
